# rootfs/etc/cron.d/rotate.cron
*/6 * * * * /bin/bash /root/rotate_ip.sh
@reboot /bin/bash /root/rotate_ip.sh

